# autonfs.conf
# autonfs - mount NFS shares on fileserver, if present

description     "Mount NFS-shares"

start on (filesystem)
respawn

exec /usr/local/bin/autonfs.sh
